# templates/scripts/reposync.cron.epp
47 22 * * * root /root/cron/reposync > /var/log/reposync.log 2>&1
 